# system_files/etc/systemd/system/24hg-heartbeat.timer
[Unit]
Description=24HG Heartbeat (boot + every 6 hours)

[Timer]
OnBootSec=3min
OnUnitActiveSec=6h
RandomizedDelaySec=30min
Persistent=true

[Install]
WantedBy=timers.target
