# system_files/etc/systemd/system/forge-parental.service
[Unit]
Description=24HG Forge Parental Controls Daemon
Documentation=man:forge-parental(1)
After=network.target graphical.target
Wants=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/forge-parental daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
PIDFile=/run/forge-parental.pid

# Security hardening
ProtectSystem=strict
ReadWritePaths=/etc/forge /var/lib/forge /var/log/forge /run
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=false
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=forge-parental

[Install]
WantedBy=multi-user.target
