# system_files/etc/systemd/system/forge-first-boot-setup.service
[Unit]
Description=24HG Forge First Boot Setup
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/forge/.first-boot-done

[Service]
Type=oneshot
ExecStart=/usr/lib/forge/first-boot-setup.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
