# system_files/etc/systemd/system/forge-auto-update.timer
[Unit]
Description=24HG Forge Update Check (boot + daily)

[Timer]
# Check on every boot (after 2 min delay for network)
OnBootSec=2min
# Also check daily at 4am
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
