# system_files/etc/systemd/system/hubos-auto-update.service
[Unit]
Description=HubOS Automatic Update Check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/lib/hubos/auto-update.sh
