# system_files/etc/systemd/user/24hg-cleaner.timer
[Unit]
Description=24HG Cleaner — Weekly cleanup timer

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
