# system_files/etc/systemd/user/24hg-cleaner.service
[Unit]
Description=24HG Cleaner — Clean gaming junk files

[Service]
Type=oneshot
ExecStart=/usr/bin/24hg-cleaner clean
