# system_files/etc/systemd/user/24hg-perks-claim.service
[Unit]
Description=24HG Perks Auto-Claim
After=24hg-hub-bridge.service
Wants=24hg-hub-bridge.service

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 30
ExecStart=/usr/bin/24hg-perks claim --auto
RemainAfterExit=yes

[Install]
WantedBy=default.target
