# system_files/etc/systemd/user/24hg-backup.service
[Unit]
Description=24HG Game Save Auto-Backup

[Service]
Type=oneshot
ExecStart=/usr/bin/24hg-backup auto
