# system_files/etc/systemd/user/24hg-wallpaper.service
[Unit]
Description=24HG Wallpaper Rotation

[Service]
Type=oneshot
ExecStart=/usr/bin/24hg-wallpaper rotate
