# system_files/etc/systemd/user/forge-deals.service
[Unit]
Description=24HG Forge Deals Checker
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/forge-deals check
Environment=HOME=%h
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
