# system_files/etc/systemd/user/forge-digest.service
[Unit]
Description=24HG Forge Weekly Digest Generator
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/forge-digest send
Environment=HOME=%h
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
