# system_files/etc/systemd/user/forge-power-plan.service
[Unit]
Description=24HG Forge Power Plan — Auto-switch power modes on AC/battery
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/forge-power-plan daemon
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
