# system_files/etc/systemd/user/forge-game-timer.service
[Unit]
Description=24HG Forge Game Timer — Track playtime across all launchers
Documentation=https://os.24hgaming.com
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/forge-game-timer daemon
Restart=on-failure
RestartSec=30
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
