# system_files/etc/systemd/user/forge-hw-scout.service
[Unit]
Description=24HG Forge Hardware Scout — Continuous hardware monitoring
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/forge-hw-scout daemon
Restart=on-failure
RestartSec=60

[Install]
WantedBy=default.target
