# system_files/etc/systemd/user/forge-peripherals.service
[Unit]
Description=24HG Forge Peripherals — Auto-switch device profiles on game launch
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/forge-peripherals daemon
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
