# system_files/etc/systemd/user/forge-cleaner.timer
[Unit]
Description=24HG Forge Cleaner — Weekly cleanup timer

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
