# system_files/etc/systemd/user/forge-challenges.timer
[Unit]
Description=24HG Forge Challenge Tracker — hourly progress check

[Timer]
OnCalendar=hourly
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
