# system_files/etc/systemd/user/forge-achievements.service
[Unit]
Description=24HG Forge Achievement Tracker
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/forge-achievements watch
Restart=on-failure
RestartSec=60

[Install]
WantedBy=default.target
