# system_files/etc/systemd/user/forge-crash-report.service
[Unit]
Description=24HG Forge Crash Report Daemon
Documentation=https://24hgaming.com/forge/docs/crash-reporting
After=default.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/forge-crash-report daemon
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=3

# Environment
Environment=FORGE_DEBUG=0

# Resource limits — lightweight daemon
MemoryHigh=64M
MemoryMax=128M
CPUQuota=5%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/forge %h/.cache/forge %h/.config/forge
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=forge-crash-report

[Install]
WantedBy=default.target
