# system_files/etc/systemd/user/forge-wallpaper.service
[Unit]
Description=24HG Forge Wallpaper Rotation

[Service]
Type=oneshot
ExecStart=/usr/bin/forge-wallpaper rotate
