# system_files/etc/systemd/user/forge-notify.service
[Unit]
Description=24HG Forge Notification Center — Rich KDE notifications for all forge tools
Documentation=https://os.24hgaming.com
After=graphical-session.target dbus.socket
Requires=dbus.socket

[Service]
Type=simple
ExecStartPre=/usr/bin/mkdir -p /tmp/forge-notifications
ExecStart=/usr/bin/forge-notify daemon
ExecStop=/usr/bin/forge-notify stop
Restart=on-failure
RestartSec=10
Environment=DISPLAY=:0

[Install]
WantedBy=graphical-session.target
