# system_files/etc/systemd/user/forge-backup.service
[Unit]
Description=24HG Forge Game Save Auto-Backup

[Service]
Type=oneshot
ExecStart=/usr/bin/forge-backup auto
