# system_files/etc/systemd/user/forge-quick-resume.service
[Unit]
Description=24HG Forge Quick Resume — Auto-restore session on login
After=graphical-session.target
ConditionPathExists=%h/.config/forge/quick-resume/auto-restore

[Service]
Type=oneshot
ExecStart=/usr/bin/forge-quick-resume auto-restore
RemainAfterExit=yes
