# system_files/etc/systemd/user/forge-achievements.service
[Unit]
Description=24HG Forge Achievement Check

[Service]
Type=oneshot
ExecStart=/usr/bin/forge-achievements check
