# system_files/etc/systemd/user/forge-cleaner.service
[Unit]
Description=24HG Forge Cleaner — Clean gaming junk files

[Service]
Type=oneshot
ExecStart=/usr/bin/forge-cleaner clean
