# system_files/etc/systemd/user/hubos-backup.timer
[Unit]
Description=HubOS Weekly Game Save Backup

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
