# system_files/etc/systemd/user/hubos-achievements.service
[Unit]
Description=HubOS Achievement Check

[Service]
Type=oneshot
ExecStart=/usr/bin/hubos-achievements check
