# system_files/etc/systemd/user/hubos-smart-launch.service
[Unit]
Description=HubOS Smart Game Launch Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/hubos-smart-launch daemon
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
