# system_files/etc/systemd/user/hubos-perks-claim.service
[Unit]
Description=HubOS Perks Auto-Claim
After=hubos-hub-bridge.service
Wants=hubos-hub-bridge.service

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 30
ExecStart=/usr/bin/hubos-perks claim --auto
RemainAfterExit=yes

[Install]
WantedBy=default.target
