# system_files/etc/systemd/user/hubos-discord-fix.service
[Unit]
Description=HubOS Discord Rich Presence Bridge
Documentation=https://os.24hgaming.com
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/hubos-discord-fix monitor
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
