# system_files/etc/systemd/user/hubos-anticheat-tracker.service
[Unit]
Description=HubOS Anti-Cheat Status Tracker
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/hubos-anticheat-tracker watch
