# system_files/etc/systemd/user/hubos-wallpaper.service
[Unit]
Description=HubOS Wallpaper Rotation

[Service]
Type=oneshot
ExecStart=/usr/bin/hubos-wallpaper rotate
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
