# system_files/etc/systemd/user/hubos-hub-bridge.service
[Unit]
Description=HubOS Hub Bridge — 24HG notifications, DMs, friend status
Documentation=https://os.24hgaming.com
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/hubos-hub-bridge daemon
Restart=on-failure
RestartSec=30
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
