# system_files/etc/systemd/user/hubos-backup.service
[Unit]
Description=HubOS Game Save Auto-Backup

[Service]
Type=oneshot
ExecStart=/usr/bin/hubos-backup auto
